# conky/.config/autostart/conky.desktop
[Desktop Entry]
Exec=conky -d -p 10
Hidden=true
Name=conky
StartupNotify=false
Terminal=false
Type=Application
