# conky/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=conky
Exec=conky -d -p 10
StartupNotify=false
Terminal=false
